# bbs.conf
[program:flask_bbs]
command=/usr/local/bin/gunicorn wsgi --bind localhost:2000
directory=/var/www/flask_bbs
autostart=true
autorestart=true

stdout_logfile=/var/log/supervisor/bbs_stdout.log
stderr_logfile=/var/log/supervisor/bbs_stderr.log
